# deploy/macos/com.teaseme.staging-web.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.teaseme.staging-web</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/zsh</string>
      <string>/Users/runner/tease-me-staging/scripts/serve-dist.sh</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
      <key>PORT</key>
      <string>4173</string>
      <key>HOST</key>
      <string>0.0.0.0</string>
      <key>DIST_DIR</key>
      <string>/Users/runner/tease-me-staging/dist</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/runner/tease-me-staging</string>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/runner/Library/Logs/tease-me-staging-web.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/runner/Library/Logs/tease-me-staging-web.log</string>
  </dict>
</plist>
